# rootfs/bootscripts.magic1/etc/profile.d/i18n.sh
# Set up i18n variables
export LANG="zh_CN.GB18030"

