# rootfs/bootscripts/etc/profile.d/dircolors.sh
# Setup for /bin/ls to support color, the alias is in /etc/bashrc.
alias ls='ls --color=auto'

