# rootfs/bootscripts.mgc30_64/usr/lib/systemd/system/uuidd.socket
[Unit]
Description=UUID daemon activation socket

[Socket]
ListenStream=/var/uuidd/request

[Install]
WantedBy=sockets.target
