# buildpkg/bootscripts/etc/udev/rules.d/90-alsa.rules
SUBSYSTEM=="sound", KERNEL=="controlC*"	RUN+="/sbin/salsa"
SUBSYSTEM=="sound", KERNEL=="pcm*"	RUN+="/sbin/salsa"		