# rootfs/bootscripts.arch1/etc/profile.d/umask.sh
# By default we want the umask to get set.
umask 022
